# etc/crontabs/root
*/1  * * * * /tmp/usr/bin/bismark-measure-wrapper
*/1  * * * * /usr/bin/bismark-probe
*/1  * * * * /usr/bin/paws-probe
*/5  * * * * /usr/bin/bismark-rshell
*/5  * * * * /usr/bin/paws-rshell
*/5  * * * * /usr/bin/bismark-sshd
*/5  * * * * ash /etc/hotplug.d/iface/20-ntpclient ifup
0 */12 * * * /usr/bin/bismark-action scriptupdate
0 5,17 * * * /usr/bin/bismark-experiments-manager -d
0 5,17 * * * /usr/bin/bismark-updater -d
5 */12 * * * /usr/bin/bismark-action mgmtconfupdate
